# roles/cron/files/crontab
*/15 * * * * root /usr/local/bin/github-bot
*/15 * * * * root /usr/local/bin/github-sucker
*/15 * * * * root /usr/local/bin/vm-admin dump
